# docs/file-warden.example.service
[Unit]
Description=File monitoring service

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target

